# homeassistant/addon/rootfs/etc/services.d/printernizer/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Printernizer Service Finish Handler
# Handles cleanup when the service stops
# ==============================================================================

bashio::log.info "Printernizer service stopped"

# If the service crashes, wait before restart
if [ "${1}" -ne 0 ] && [ "${1}" -ne 256 ]; then
    bashio::log.warning "Printernizer crashed with exit code ${1}"
    bashio::log.info "Waiting 10 seconds before restart..."
    sleep 10
fi